# target/linux/airoha/an7581/base-files/etc/uci-defaults/04_gemtekfan
#
# Copyright (C) 2025 openwrt.org
#

. /lib/functions.sh

board=$(board_name)

case "$board" in
gemtek,w1700k)
	# Disable automatic control temporarily to safely write
	echo 1 > /sys/class/hwmon/hwmon5/pwm1_enable

	# Apply new PWM values
	echo 54  > /sys/class/hwmon/hwmon5/pwm1_auto_point1_pwm
	echo 69  > /sys/class/hwmon/hwmon5/pwm1_auto_point2_pwm
	echo 95  > /sys/class/hwmon/hwmon5/pwm1_auto_point3_pwm
	echo 199 > /sys/class/hwmon/hwmon5/pwm1_auto_point4_pwm

	# Re-enable automatic mode
	echo 2 > /sys/class/hwmon/hwmon5/pwm1_enable
	;;
esac

exit 0
